# deploy/systemd/bi-rtx-api.service
[Unit]
Description=BI-IDE RTX AI API Server (Port 8090)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bi
WorkingDirectory=/home/bi/.bi-ide-worker
ExecStart=/home/bi/.bi-ide-worker/venv/bin/python3 /home/bi/.bi-ide-worker/rtx4090_machine/rtx_api_server.py
Restart=always
RestartSec=5
Environment=HOME=/home/bi
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/bi/.bi-ide-worker

[Install]
WantedBy=multi-user.target
